# relay-server/prexu-relay.service
[Unit]
Description=Prexu Watch Together Relay Server
After=network.target

[Service]
Type=simple
User=plex
ExecStart=/usr/local/bin/prexu-relay --port 8080
Restart=always
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
